# frontend/download/deploy/nextjs-lorcgr.service
[Unit]
Description=LOR-CGR Next.js Frontend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/lorcgr/frontend
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0

[Install]
WantedBy=multi-user.target
